# config/systemd/ak-historian-r1p5r-readiness-watch.timer
[Unit]
Description=Run AK Historian R1P5R structural readiness watch hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=0
Unit=ak-historian-r1p5r-readiness-watch.service

[Install]
WantedBy=timers.target
